# ops/systemd/genapi-shim-chat-responses.service.example
[Unit]
Description=GenAPI generic chat->responses shim
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/.ductor/workspace/web_capability_api/shims/chat-responses
EnvironmentFile=-/root/.ductor/workspace/web_capability_api/ops/env/providers/chat-responses.env
ExecStart=/usr/local/bin/node /root/.ductor/workspace/web_capability_api/shims/chat-responses/chat_responses_shim.mjs
Restart=on-failure
RestartSec=3
LimitNOFILE=100000

[Install]
WantedBy=multi-user.target
